# scoreboard/bin/upstart/ictf-gunicorn.conf
description "Gunicorn application server handling ictf API"

start on runlevel [2345]
stop on runlevel [!2345]

respawn
#setuid deploy
#setgid deploy
chdir /opt/ictf/scoreboard

exec /opt/ictf/scoreboard/bin/run_app.sh
