# dashboard/dashboard.conf
description "The CTF web interface"

# When to start the service
start on started ctf-database

# When to stop the service
stop on runlevel [016]

# Automatically restart process if crashed
respawn
respawn limit unlimited

console log

chdir /opt/dashboard

pre-start script
    exec iptables -t nat -A PREROUTING -p tcp --dport 80 -j REDIRECT --to 4567
end script

# Start the process
script
    exec bundle exec foreman start
end script
