# router/services/ictf-pcap-s3.conf
description "iCTF tcpdump Service"

# When to start the service
start on runlevel [2345]

# When to stop the service
stop on runlevel [016]

# Automatically restart process if crashed
respawn
respawn limit unlimited

console log

#setuid gamemaster

chdir /ictf/router

# Start the process
exec /usr/bin/python ictf-pcap-s3.py
